# docker-allinone/supervisord.conf
[supervisord]
nodaemon=true

[program:webapp]
command=/opt/bin/start-webapp.sh
autostart=true # 自动启动
autorestart=false
startsecs=0
startretries=0

;Logs (all Hub activity redirected to stdout so it can be seen through "docker logs"
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:selenoid]
command=/usr/bin/selenoid -listen :4444 -conf /etc/selenoid/browsers.json -limit 100 -container-network bridge
autostart=true
autorestart=false
startsecs=0
startretries=0

;Logs (all Hub activity redirected to stdout so it can be seen through "docker logs"
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:qqbot]
command=/go-cqhttp/go-cqhttp -c /go-cqhttp/config.yml -d
autostart=true
autorestart=false
startsecs=0
startretries=0

;Logs (all Hub activity redirected to stdout so it can be seen through "docker logs"
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0